# overlay/root/bin/install_zj2m
#!/usr/bin/env bash

# shellcheck disable=SC1091
. /etc/rc.subr && load_rc_config zwavejs2mqtt
: "${zwavejs2mqtt_app_dir:=/usr/local/share/zwavejs2mqtt}"

workdir="$(mktemp -dt build_src)"

get_latest_release() {
  local tarball_url tarball
  tarball_url="$(curl -s https://api.github.com/repos/zwave-js/zwavejs2mqtt/releases/latest | grep tarball_url | cut -d : -f 2,3 | tr -d \",)"
  tarball="${workdir}/source.tar.gz"
  echo -ne "\nDownloading latest release..."
  # shellcheck disable=SC2086
  wget --quiet -O "${tarball}" ${tarball_url}
  tar -xzf "${tarball}" -C "${workdir}"
  cd "${workdir}"/*zwavejs2mqtt* || return 1
  echo "  OK"
}

build_zj2m() {
  ## This build function is based on zwavejs2mqtt/docker/Dockerfile
  local updateDevices
  local zwavejs="https://github.com/zwave-js/node-zwave-js"
  npm config set unsafe-perm true
  npm install
  npm run build
  npm prune --production
  rm -rf -- \
      build \
      docker \
      docs \
      index.html \
      kubernetes \
      kustomization.yaml \
      package-lock.json \
      package.sh \
      src \
      static\
      store \
      stylesheets
  if [ -n "${updateDevices}" ]; then
    git clone "${zwavejs}"
    rsync -av node-zwave-js/packages/config/config/devices/ './node_modules/@zwave-js/config/config/devices'
    rm -rf -- node-zwave-js
  fi
}

copy_zj2m() {
  local app_dir="${zwavejs2mqtt_app_dir}"
  ## remove all (dot) .files and .folders from zj2m build_src directory
  find . -mindepth 1 -name '.*' -exec rm -rf -- "{}" +
  ## if it's enabled and running, stop the zwavejs2mqtt service and set it to restart
  if sysrc zwavejs2mqtt_enable &>/dev/null && service zwavejs2mqtt status &>/dev/null; then
    local restart='YES'
    service zwavejs2mqtt stop || err 1 "failed to stop service"
  fi
  ## if present, remove existing app directory
  if [ -d "${app_dir}" ]; then
    rm -rf -- "${app_dir}"
  fi
  ## install new app directory and copy zj2m build_src to it
  install -d -m 775 -- "${app_dir}" \
  && cp -Rf . "${app_dir}"
  ## restart the service, if it was previously set to restart
  if [ "${restart}" == 'YES' ]; then
    service zwavejs2mqtt start
  fi
}

get_latest_release \
&& build_zj2m \
&& copy_zj2m \
|| echo -e "\ninstall failed :("
rm -rf -- "${workdir}"
